# crw/rootfs/etc/s6-overlay/s6-rc.d/searxng/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

mkdir -p /run/crw

# botdetection probes for limiter.toml next to the settings file even with
# the limiter off — an empty file silences the startup warning.
touch /run/crw/limiter.toml

# Per-boot secret: the JSON API is stateless, nothing needs to survive restarts.
SECRET_KEY="$(head -c 32 /dev/urandom | od -A n -t x1 | tr -d ' \n')"

/opt/venv-searxng/bin/python /opt/mcp-bridge/render_settings.py \
  --base /etc/searxng/settings.yml \
  --options /data/options.json \
  --secret "${SECRET_KEY}" \
  --out /run/crw/searxng-settings.yml

export SEARXNG_SETTINGS_PATH=/run/crw/searxng-settings.yml

echo "[searxng] starting on 127.0.0.1:8080 …"
exec /opt/venv-searxng/bin/granian --interface wsgi \
  --blocking-threads 4 \
  --host 127.0.0.1 --port 8080 searx.webapp:app
