# crw/rootfs/etc/s6-overlay/s6-rc.d/crw/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

# fastcrw sidecar pattern: search traffic flows exclusively through
# /v1/search, backed by the bundled SearXNG on loopback.
export CRW_SERVER__HOST=127.0.0.1
export CRW_SERVER__PORT=3000
export CRW_SEARCH__SEARXNG_URL=http://127.0.0.1:8080

# Optional egress proxy for scraping (anti-bot bypass). reqwest picks up the
# standard proxy env vars; loopback services must stay direct.
PROXY="$(bashio::config 'outgoing_proxy' '' 2>/dev/null || true)"
if [ -n "${PROXY}" ] && [ "${PROXY}" != "null" ]; then
  export HTTP_PROXY="${PROXY}" HTTPS_PROXY="${PROXY}"
  export NO_PROXY="127.0.0.1,localhost"
  echo "[crw] outgoing proxy enabled"
fi

# s6 only orders process start, not readiness — wait for SearXNG to listen
# so crw doesn't log its search backend as unreachable on every boot.
for _ in $(seq 1 60); do
  if curl -sf -m 1 -o /dev/null -H 'X-Forwarded-For: 127.0.0.1' http://127.0.0.1:8080/; then
    echo "[crw] searxng is ready"
    break
  fi
  sleep 0.5
done

echo "[crw] starting crw-server on 127.0.0.1:3000 …"
exec /usr/local/bin/crw-server
