# crw/rootfs/etc/s6-overlay/s6-rc.d/crw/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -euo pipefail

# fastcrw sidecar pattern: search traffic flows exclusively through
# /v1/search, backed by the bundled SearXNG on loopback.
export CRW_SERVER__HOST=127.0.0.1
export CRW_SERVER__PORT=3000
export CRW_SEARCH__SEARXNG_URL=http://127.0.0.1:8080

echo "[crw] starting crw-server on 127.0.0.1:3000 …"
exec /usr/local/bin/crw-server
